# agent/packaging/rmmagent.service
[Unit]
Description=rmmagic endpoint agent
Documentation=https://github.com/codex666-cenotaph/rmmagic
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/rmmagent run --state-dir /var/lib/rmmagent
Restart=always
RestartSec=5s

# Hardening. The agent executes scripts and installs OS packages on the
# endpoint, so it runs as root with broad privilege — keep these settings
# conservative and tighten per your threat model.
NoNewPrivileges=false
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
# /var/lib/rmmagent holds identity/journal/update state. /usr/bin is also
# writable so signed auto-update can atomically swap the agent binary in
# place (kept narrow to the install dir; ProtectSystem=full still guards the
# rest of /usr, /boot, /etc). Restart=always above relaunches into the new
# binary after a swap; the agent's watchdog rolls back on failed health check.
ReadWritePaths=/var/lib/rmmagent /usr/bin

[Install]
WantedBy=multi-user.target
